# crontab
0 */6 * * * /srv/ocs/update-software-inventory.sh
