# crontab
0 */6   * * *   root    /srv/ocs/update-software-inventory.sh
